# automation/docker/conf/supervisord.conf
[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[unix_http_server]
file = /run/supervisord.sock

[program:uvicorn]
command=/usr/local/bin/uvicorn main:app --port 80 --host 0.0.0.0 --workers %(ENV_WORKERS)s
process_name=uvicorn
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=8
startsecs=4
autorestart=false
